# package/ezp-power-detect/files/power-detect.init
#!/bin/sh /etc/rc.common
BIN=/usr/sbin/power-detect

START=98

log_msg="Power-detect-server"
start () {
    $BIN || err=1 
    return $err
}
