# package/ezp-i2c/files/ezp-i2c-init
#!/bin/sh /etc/rc.common
. /etc/arch.sh

START=09
log_msg="I2C"
boot() {
    mknod /dev/i2c c 218 0
}
