# package/base-files/files/etc/hotplug.d/net/11-dhcp
#!/bin/sh
./etc/arch.sh

[ "$(nvram get wan0_proto)" = "htc" -o "$(nvram get wan0_proto)" = "iphone" ] && {
killall udhcpc
}


