# package/ezp-eth-led/files/ezp-eth-led.init
#!/bin/sh /etc/rc.common
BIN=/usr/sbin/ezp-eth-led

START=98

log_msg="EZP-eth-server"
start () {
    $BIN || err=1 
    return $err
}
